# addon/harbor-companion/rootfs/etc/s6-overlay/s6-rc.d/harbor-companion/run
#!/command/execlineb -P
python3 /app/main.py
